# distro/apline/3.20/cacti-lighttpd/docs/cactid
#!/sbin/openrc-run

name="cacti"
command="php /home/cacti/cactid.php"
command_background=true
pidfile="/var/run/${RC_SVCNAME}.pid"
command_args=""
command_user="cacti"
command_group="lighttpd"
output_log="/var/log/cacti/${RC_SVCNAME}.log"
error_log="/var/log/cacti/${RC_SVCNAME}.err"
depend() {
    need net
}

start_pre() {
    checkpath --directory --owner cacti:cacti --mode 0755 /var/log/cacti
}